# profile.d/expidus-sdk.sh.in
if [[ ! -z "$EXPIDUS_SDK_VERSION" ]]; then
  source "@LIBDIR@/expidus-sdk/$EXPIDUS_SDK_VERSION/ensure.sh"
else
  if [[ -e @LIBDIR@/expidus-sdk/latest ]]; then
    source @LIBDIR@/expidus-sdk/latest/@TARGET_SYSTEM@/ensure.sh
  else
    source @SDK_LIBDIR@/ensure.sh
  fi
fi
