# profile.d/meson.build
configure_file(input: 'expidus-sdk.sh.in',
  output: 'expidus-sdk.sh',
  configuration: conf_data,
  install: true,
  install_dir: join_paths(sysconfdir, 'profile.d'),
  install_mode: 'rwxr-xr-x')
